# app-misc/alfresco-bin/files/alfresco.init
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

source /usr/share/tomcat-scripts/runscript.sh

description="Alfresco, Open Source Enterprise Content Management System"

catalina_home=
catalina_base=
catalina_temp=

tomcat_extra_jars=

depend() {
	use net
}
