# www-apps/xwiki-enterprise-bin/files/xwiki.init
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

source /usr/share/tomcat-scripts/runscript.sh

description="XWiki Enterprise, an extensible wiki platform writtern in Java."

catalina_home=
catalina_base=
catalina_temp=

tomcat_extra_jars=

depend() {
	use net
}

